# config/set-pcie-latency.service
[Unit]
Description=Set PCIe Latency Timers
After=multi-user.target
Wants=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/set-pcie-latency.sh
RemainAfterExit=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

